# nvidia/system/nvidia-powerd.service
[Unit]
Description=nvidia-powerd service

[Service]
Type=simple
ExecStart=/usr/bin/nvidia-powerd

[Install]
WantedBy=multi-user.target
